# Ansible/roles/Configure_Proxmox/files/zfs-scrub@.timer
[Unit]
Description=Scrub a ZFS Pool

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
